# rtorrentfloodvpn/build/flood.conf
[program:flood]
autorestart = true
startsecs = 5
user = nobody
command = /usr/bin/node --use_strict /usr/lib/flood/index.js --rundir /config/flood --host 0.0.0.0 --port 3000 --allowedpath /config --allowedpath /downloads --dehost delugevpn --deport 8112 --deuser admin --depass 8Drag0n4 --rthost 127.0.0.1 --rtport 5000
umask = 000
redirect_stderr=true
redirect_stdout=true
